# installation/derban.service
[Unit]
Description=Derban

[Service]
Type=simple
ExecStart=/usr/bin/derban.sh

[Install]
WantedBy=multi-user.target